# crontab
PATH=/go/bin:/usr/local/go/:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/30    *   * * *	root	/run.sh >>/var/log/ptndown.log
